# focusblock-disable.timer
[Unit]
Description=Disable focusblock daily at 17:00

[Timer]
OnCalendar=*-*-* 17:00:00
Persistent=true

[Install]
WantedBy=timers.target
